# site-cookbooks/deploy/templates/capybara-webkit-daemon.service.erb
[Unit]

Description=Capybara Webkit daemon

After=network.target
Requires=xvfb.service



[Service]

Type=simple
Restart=always

User=<%= @user %>
Group=<%= @group %>

WorkingDirectory=/srv/scrapod/current/
PIDFile=/srv/scrapod/shared/tmp/pids/capybara.pid

Environment=BUNDLE_GEMFILE=/srv/scrapod/current/Gemfile

<% if @display.nil? %>
ExecStart=/usr/bin/bundle exec --keep-file-descriptors capybara-webkit-daemon --debug --logfile /srv/scrapod/shared/log/capybara.log
<% else %>
ExecStart=/usr/bin/bundle exec --keep-file-descriptors capybara-webkit-daemon --debug --logfile /srv/scrapod/shared/log/capybara.log --display <%= @display %>
<% end %>

ExecStop=/bin/kill `cat /srv/scrapod/shared/tmp/pids/capybara.pid`
