# site-cookbooks/scrapod/templates/capybara.service.erb
<% bundle = '/usr/bin/bundle exec --keep-file-descriptors' %>

[Unit]
Description=Capybara Webkit daemon
After=network.target
Requires=xvfb.service

[Service]
Type=simple
Restart=always
User=<%= @user %>
Group=<%= @group %>

WorkingDirectory=/srv/scrapod/current/
PIDFile=/srv/scrapod/shared/tmp/pids/capybara.pid
Environment=BUNDLE_GEMFILE=/srv/scrapod/current/Gemfile

<% if @display.nil? %>
ExecStart=<%= bundle %> capybara-webkit-daemon --config /srv/scrapod/current/config/capybara_webkit_daemon/production.rb
<% else %>
ExecStart=<%= bundle %> capybara-webkit-daemon --config /srv/scrapod/current/config/capybara_webkit_daemon/production.rb --display <%= @display %>
<% end %>

[Install]
WantedBy=multi-user.target
